# ansible/roles/loki/templates/loki.service.j2
[Unit]
Description=Loki Log Aggregation System
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ lgtm_user }}
Group={{ lgtm_group }}
ExecStart={{ lgtm_base_dir }}/bin/loki \
  -config.file={{ loki_config_dir }}/loki-config.yml \
  -target=all

ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=loki

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ loki_data_dir }} {{ loki_config_dir }}

# Resource limits
LimitNOFILE={{ nofile_limit }}
LimitNPROC={{ nproc_limit }}

[Install]
WantedBy=multi-user.target
